# crontab
0 0 * * * PATH='/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin' apt-get update && apt-get dist-upgrade --yes && apt-get autoremove --purge --yes
