# foobar.service
[Unit]
Description=Just put me on /lib/systemd/system/ or $HOME/.config/systemd/user/, then systemctl enable <it>. If you want to keep application running, change Type to forking
After=network.target

[Service]
Type=simple
User=nobody
Restart=on-failure
RestartSec=5s
WorkingDirectory=/
ExecStart=/bin/bash -c "echo foobar"

[Install]
WantedBy=multi-user.target
