# foobar.service
[Unit]
Description=Foobar Service, just put me on /lib/systemd/system/
After=network.target

[Service]
Type=simple
User=nobody
Restart=on-failure
RestartSec=5s
WorkingDirectory=/
ExecStart=/bin/bash -c "echo foobar"

[Install]
WantedBy=multi-user.target
