# supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile = /tmp/supervisord.log
nodaemon=true

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:hello]
command=bash -c "echo 'Hello, World!'"
autostart=true
startsecs=0
autorestart=unexpected
exitcodes=0,2
directory=/tmp
environment=A="1",B="2"
user=root
