# supervisor-rc-local.ini
[program:username-rc-local]
command=/opt/homebrew/bin/bash -c '/Users/username/etc/rc.local'
user=username
autorestart=false
startsecs=0
environment=HOME="/Users/username",USER="username",PATH="/opt/homebrew/bin:/usr/bin:/bin:%(ENV_PATH)s"
